# deploy/supervisor/psych-web-scheduler.conf
; Laravel scheduler for psych-web. Runs `schedule:work` as a long-lived
; process instead of a cron entry, so it's managed the same way as the
; queue worker and Reverb. Drives `rooms:prune` (routes/console.php) every
; five minutes to delete empty/inactive rooms.
;
; Install:
;   cp deploy/supervisor/psych-web-scheduler.conf /etc/supervisor/conf.d/
;   # edit the path/user placeholders below first
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start psych-web-scheduler:*

[program:psych-web-scheduler]
command=php /var/www/psych-web/artisan schedule:work
directory=/var/www/psych-web
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/www/psych-web/storage/logs/supervisor-scheduler.log
stdout_logfile_maxbytes=10MB
