# deploy/supervisor/psych-web-queue.conf
; Laravel queue worker for psych-web (QUEUE_CONNECTION=database).
; Drives ExpireAnswerPhase / ExpireVotingPhase / ExpireRevealPhase — without
; this running, rounds never auto-advance once their timer runs out.
;
; Install:
;   cp deploy/supervisor/psych-web-queue.conf /etc/supervisor/conf.d/
;   # edit the path/user placeholders below first
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start psych-web-queue:*

[program:psych-web-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/psych-web/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --timeout=90
directory=/var/www/psych-web
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
redirect_stderr=true
stdout_logfile=/var/www/psych-web/storage/logs/supervisor-queue.log
stdout_logfile_maxbytes=10MB
