# deploy/supervisor/psych-web-reverb.conf
; Laravel Reverb WebSocket server for psych-web (real-time game broadcasts).
; Match --host/--port to REVERB_HOST / REVERB_PORT in .env (defaults below
; assume the app is proxied, so Reverb itself listens on all interfaces).
;
; Install:
;   cp deploy/supervisor/psych-web-reverb.conf /etc/supervisor/conf.d/
;   # edit the path/user placeholders below first
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start psych-web-reverb:*

[program:psych-web-reverb]
command=php /var/www/psych-web/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/psych-web
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/www/psych-web/storage/logs/supervisor-reverb.log
stdout_logfile_maxbytes=10MB
